# ops/systemd/agenttrader-market-ingest.service
[Unit]
Description=AgentTrader VM Market Ingest
Wants=network-online.target
After=network-online.target

# Prevent tight restart loops when config/code isn't present.
ConditionPathExists=/etc/agenttrader/market-ingest.env
ConditionPathExists=/opt/agenttrader/agent-trader-v2

[Service]
Type=simple

# Run as a dedicated, non-root user.
User=agenttrader
Group=agenttrader

WorkingDirectory=/opt/agenttrader/agent-trader-v2

# Inject secrets/config without printing them.
EnvironmentFile=/etc/agenttrader/market-ingest.env

# Make logs immediately visible and avoid writing __pycache__ to disk.
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONPATH=/opt/agenttrader/agent-trader-v2

ExecStart=/usr/bin/python3 -m backend.ingestion.market_data_ingest

# Restart hardening: backoff + circuit-breaker.
Restart=on-failure
RestartSec=30s
StartLimitIntervalSec=10min
StartLimitBurst=3

# Clean shutdown (lets Python handle SIGINT consistently).
KillSignal=SIGINT
TimeoutStopSec=30s

# Journald best-practice (queryable by unit + identifier).
SyslogIdentifier=agenttrader-market-ingest
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

